# quectel-app/setwifimac/setwifimac.service.in
[Unit]
Description=Set WIFI MAC
DefaultDependencies=no
#After=btproperty.service

[Service]
RemainAfterExit=yes
Restart=no
#ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/setwifimac

[Install]
WantedBy=multi-user.target
