# quectel-app/ql-nv-settings/ql_nv_settings.service.in
[Unit]
Description=Default NV settings service

[Service]
Type=oneshot
#ExecStartPre=/bin/sleep 6
ExecStart=/usr/bin/ql-nv-settings 1
#ExecStop=/sbin/start-stop-daemon -K -n ql-nv-settings
Restart=no
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
