# quectel-core/qlrild/lib/systemd/system/qlrild.service
[Unit]
Description=qlrild service
DefaultDependencies=no
After=qcrild2.service

[Service]
Type=simple
ExecStart=/usr/bin/qlrild
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
